# ops/pocketbase/pocketbase-backup.service
[Unit]
Description=Backup PocketBase data

[Service]
Type=oneshot
ExecStart=/opt/pocketbase/backup-pocketbase.sh
